# deploy/macos/com.romiv.hh-ollama-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.romiv.hh-ollama-tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/ssh</string>
        <string>-N</string>
        <string>-T</string>
        <string>-i</string>
        <string>/Users/alfa/.ssh/hh_ollama_tunnel</string>
        <string>-o</string>
        <string>IdentitiesOnly=yes</string>
        <string>-o</string>
        <string>BatchMode=yes</string>
        <string>-o</string>
        <string>ExitOnForwardFailure=yes</string>
        <string>-o</string>
        <string>ConnectTimeout=15</string>
        <string>-o</string>
        <string>ServerAliveInterval=30</string>
        <string>-o</string>
        <string>ServerAliveCountMax=3</string>
        <string>-R</string>
        <string>127.0.0.1:11434:127.0.0.1:11434</string>
        <string>hh-ollama-tunnel@213.176.115.241</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/alfa/Library/Logs/hh-ollama-tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alfa/Library/Logs/hh-ollama-tunnel.log</string>
</dict>
</plist>
